# scripts/infra/supabase.cron
0 */6 * * * /bin/bash /app/scripts/supabase-ping.sh